# deploy/com.m-hero-vip-custom-alert.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.m-hero-vip-custom-alert.watchdog</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/i/myCode/m-hero-vip-custom-alert/.venv/bin/python</string>
		<string>/Users/i/myCode/m-hero-vip-custom-alert/scripts/health_watchdog.py</string>
		<string>--url</string>
		<string>http://127.0.0.1:9002/api/vip/status</string>
		<string>--interval</string>
		<string>60</string>
		<string>--cooldown</string>
		<string>1800</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/i/myCode/m-hero-vip-custom-alert</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>ThrottleInterval</key>
	<integer>5</integer>
	<key>StandardOutPath</key>
	<string>/Users/i/Library/Logs/com.m-hero-vip-custom-alert.watchdog.out.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/i/Library/Logs/com.m-hero-vip-custom-alert.watchdog.err.log</string>
</dict>
</plist>
