# deploy/com.m-hero-vip-custom-alert.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.m-hero-vip-custom-alert.web</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/i/myCode/m-hero-vip-custom-alert/.venv/bin/python</string>
		<string>/Users/i/myCode/m-hero-vip-custom-alert/scripts/web_console.py</string>
		<string>--host</string>
		<string>127.0.0.1</string>
		<string>--port</string>
		<string>9002</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/i/myCode/m-hero-vip-custom-alert</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>TZ</key>
		<string>Asia/Shanghai</string>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>ThrottleInterval</key>
	<integer>5</integer>
	<key>StandardOutPath</key>
	<string>/Users/i/Library/Logs/com.m-hero-vip-custom-alert.web.out.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/i/Library/Logs/com.m-hero-vip-custom-alert.web.err.log</string>
</dict>
</plist>
